# launchagents/Library/LaunchAgents/sg.lexi.ctx-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>sg.lexi.ctx-monitor</string>

	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/python3</string>
		<string>/Users/kiran/.claude/tools/ctx-monitor/ctx-monitor.py</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/kiran/.claude/tools/ctx-monitor</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>ProcessType</key>
	<string>Background</string>

	<key>StandardOutPath</key>
	<string>/dev/null</string>

	<!-- NOT /tmp/claude-ctx: launchd opens this file BEFORE exec'ing the job and
	     does not create intermediate directories. ctx-monitor.py makes its own
	     state dir, but that runs too late, and macOS reaps untouched /tmp files
	     after ~3 days — so a missing /tmp/claude-ctx meant launchd could fail to
	     spawn the job at login. ~/Library/Logs always exists and is never reaped. -->
	<key>StandardErrorPath</key>
	<string>/Users/kiran/Library/Logs/ctx-monitor.err</string>
</dict>
</plist>
